# etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin ft --noissue --nonewline --skip-login --noclear %I $TERM
