# etc/systemd/system/gpu-reboot.service
[Unit]
Description=Reboot

[Service]
Type=simple
ExecStart=/usr/sbin/reboot now

