# etc/systemd/system/gpu-reboot.timer
[Unit]
Description=GPU bug workaround reboot timer

[Timer]
OnCalendar=*-*-* 3:00:00

[Install]
WantedBy=multi-user.target
