# home/ft/.config/systemd/user/pointcloud-server@.service
[Unit]
Description=Pointcloud server
[Service]
ExecStart=python3 -mfootron_pc_server
Restart=always
[Install]
WantedBy=X11@%i.target
